# ops/auditd-rules.conf
# ops/auditd-rules.conf — kernel-level audit rules for Mafiking.
#
# Installed at /etc/audit/rules.d/99-mafiking.rules
# Loaded with:  augenrules --load   (or restart auditd)
#
# CIS 4.1.x — syscall auditing for sensitive operations.
# ASVS V7.1, V7.2, V7.3.
#
# Notes:
#   - Excludes the 9router maintenance process to avoid noise.
#   - Buffers are tuned for 957 MB RAM; if you scale up, raise the buffers.

# --- Buffer ---
-b 8192
-f 1
--backlog_wait_time 60000

# --- Failure mode ---
-f 2

# --- Identity & auth events ---
-w /etc/passwd -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/group  -p wa -k identity
-w /etc/gshadow -p wa -k identity
-w /etc/sudoers -p wa -k sudoers
-w /etc/sudoers.d/ -p wa -k sudoers

# --- SSH ---
-w /etc/ssh/sshd_config -p wa -k sshd
-w /etc/ssh/sshd_config.d/ -p wa -k sshd

# --- App config ---
-w /opt/mafiking/.env -p wa -k mafiking-config
-w /opt/mafiking/server.js -p wa -k mafiking-config
-w /etc/nginx/sites-available/ -p wa -k nginx-config
-w /etc/modsecurity/ -p wa -k modsecurity-config

# --- Privilege escalation syscalls ---
-a always,exit -F arch=b64 -S execve -F euid=0 -F key=root-priv
-a always,exit -F arch=b64 -S setuid -F key=priv-change
-a always,exit -F arch=b64 -S setgid -F key=priv-change

# --- Network exposure (open ports / sockets) ---
-a always,exit -F arch=b64 -S bind -F key=network-bind
-a always,exit -F arch=b64 -S connect -F key=network-connect

# --- Cron ---
-w /etc/cron.d/ -p wa -k cron
-w /var/spool/cron/ -p wa -k cron

# --- Logs ---
-w /var/log/auth.log -p wa -k auth-log
-w /var/log/syslog -p wa -k syslog
-w /opt/mafiking/logs/ -p wa -k mafiking-logs

# --- Backups ---
-w /opt/mafiking-ops/ -p wa -k mafiking-backup
-w /usr/bin/rclone -p x -k backup-tool

# --- Time changes (anti-log-tampering) ---
-a always,exit -F arch=b64 -S adjtimex -S settimeofday -k time-change

# --- Module loading (anti-rootkit) ---
-w /sbin/insmod   -p x -k modules
-w /sbin/rmmod    -p x -k modules
-w /sbin/modprobe -p x -k modules

# --- Exclude 9router noise ---
-a never,user -F uid=0 -F exe=/usr/bin/9router -F key=ignored
